# pi.cron
# Display message during bootup
@reboot sudo python3 /home/pi/RocketLaunchPi/startup.py >/dev/null 2>&1

# RocketLaunchPi
5,20,35,50 08-22 * * * sudo python3 /home/pi/RocketLaunchPi/rocketlaunchpi.py >/home/pi/RocketLaunchPi/rocket.log 2>&1

